# files/systemd/tsmbackup.timer
[Unit]
Description=tsmbackup

[Timer]
OnCalendar=*-*-* 03:00:00
Unit=tsmbackup.service
Persistent=false

[Install]
WantedBy=timers.target