# files/systemd/tsmbackup.service
[Unit]
Description=tsmbackup
After=network.target

[Service]
Type=oneshot
ExecStart=/var/scripts/backup-script.sh
User=TSM_Admin
Group=TSM_Admin